# setting/docker/app-setting/supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command=uwsgi --ini /tmp/uwsgi.ini
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=10


[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
priority=20
